# ansible/roles/ru_probe/templates/reality-e2e@.timer.j2
[Unit]
Description=Run RU iperf3 probe hourly with jitter (%i)

[Timer]
OnCalendar=*-*-* %i:00:00
RandomizedDelaySec={{ ru_probe_random_delay_sec }}
Persistent=true

[Install]
WantedBy=timers.target
