# ansible/roles/ru_probe/templates/ru-iperf-probe.service.j2
[Unit]
Description=RU iperf3 throughput probe -> node_exporter textfile
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/ru-iperf-probe
ExecStart=/usr/local/bin/ru_iperf_probe.sh
User=root
Group=root
Nice=10
