# ansible/roles/ru_probe/templates/ru-iperf-probe.service.j2
[Unit]
Description=RU iperf3 throughput probe -> node_exporter textfile
After=network-online.target

[Service]
Type=oneshot
User: "{{ node_exporter_user | default('node_exporter') }}"
Group: "{{ node_exporter_group | default('node_exporter') }}"
UMask=0002
EnvironmentFile=-/etc/default/ru-iperf-probe
ExecStart=/usr/local/bin/ru_iperf_probe.sh
Nice=10
