# ansible/roles/reality_e2e/templates/reality-e2e@.timer.j2
[Unit]
Description=REALITY E2E probe timer (%i)

[Timer]
OnBootSec=30s
OnUnitActiveSec=300s
RandomizedDelaySec=60s
AccuracySec=15s
Persistent=true
Unit=reality-e2e@%i.service

[Install]
WantedBy=timers.target
